# docker/cron/crontab
* * * * * www php /app/artisan schedule:run >> /dev/stdout
